# 0.60/helium.service
[Unit]
Description=Helium Gateway

[Service]
Type=forking
ExecStart=/etc/helium_gateway/helium_gateway server
